# v1/launch_agents/com.personal-agent.web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.personal-agent.web</string>

    <key>ProgramArguments</key>
    <array>
        <string>__V1_DIR__/.venv/bin/python</string>
        <!-- -u: unbuffered stdout/stderr so log files reflect real-time
             output instead of waiting for Python's block-buffer to flush. -->
        <string>-u</string>
        <string>-m</string>
        <string>uvicorn</string>
        <string>web.app:app</string>
        <string>--host</string>
        <!-- 127.0.0.1 only — no LAN exposure, no auth boundary needed. -->
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>8770</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__V1_DIR__</string>

    <key>RunAtLoad</key>
    <true/>

    <!-- Restart on crash. ThrottleInterval prevents tight restart loops. -->
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>__V1_DIR__/data/web.log</string>
    <key>StandardErrorPath</key>
    <string>__V1_DIR__/data/web.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
